# scripts/kiosk-scripts/systemd/calendarbot-network-wait.service
[Unit]
Description=Wait for network connectivity before starting kiosk
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/local/bin/calendarbot-wait-for-network.sh
TimeoutStartSec=180
User=pi

[Install]
WantedBy=network-online.target