# kiosk/service/calendarbot-kiosk-x@.service
[Unit]
Description=CalendarBot Kiosk X Session for %i
After=network.target
# Prevent rapid restart loops
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i
# Run X server directly (no VT switching, no console requirements)
# -nolisten tcp: Don't listen on TCP (security)
# -novtswitch: Don't switch VTs
# -sharevts: Share current VT
ExecStart=/bin/sh -c 'exec /usr/bin/Xorg :0 -nolisten tcp -novtswitch -sharevts & sleep 2 && DISPLAY=:0 exec /home/%i/.xinitrc'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=DISPLAY=:0

# Grant access to devices needed for X
SupplementaryGroups=video input tty

# Grant capability to access TTY devices
AmbientCapabilities=CAP_SYS_TTY_CONFIG
CapabilityBoundingSet=CAP_SYS_TTY_CONFIG

[Install]
WantedBy=multi-user.target
