# kiosk/service/calendarbot-kiosk-x@.service
[Unit]
Description=CalendarBot Kiosk X Session for %i
After=network.target
# Prevent rapid restart loops
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i
# Start X session via startx (works because auto-login provides console user)
# Auto-login on tty1 provides the console session needed for X permissions
# This runs .xinitrc which launches the browser
ExecStart=/usr/bin/startx
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target
