# kiosk/service/calendarbot-kiosk.service
[Unit]
Description=CalendarBot Lite Calendar Display Server for %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/calendarbot
Environment="PATH=/home/%i/calendarbot/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/%i/calendarbot/venv/bin/python -m calendarbot_lite
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Resource limits (Pi Zero 2 optimized)
MemoryHigh=200M
MemoryMax=300M
CPUQuota=50%

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target