# scripts/kiosk-scripts/calendarbot-kiosk.service
[Unit]
Description=CalendarBot Python Web Server
After=network.target

[Service]
User=bencan
WorkingDirectory=/home/bencan/calendarbot
ExecStart=/bin/bash -c 'source venv/bin/activate && calendarbot --web --port 8080'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target